# config/customize/gnome/autostart/gnome-pomodoro.desktop
[Desktop Entry]
Type=Application
Name=Pomodoro
Comment=Timer
Icon=stopwatch-symbolic
Exec=gnome-pomodoro --no-default-window
